# autostart/chromium-browser.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/usr/bin/chromium-browser --no-startup-window
Name=Chromium
